# powernap/etc/systemd/system/powernapd.service
[Unit]
Description=Start powernap daemon

[Service]
Type=simple
ExecStart=/usr/local/bin/powernapd

[Install]
WantedBy=multi-user.target

